# customer/project/HC5861/target/linux/ramips/base-files/etc/uci-defaults/disable-pdnsd
#!/bin/sh
[ -x /etc/init.d/pdnsd ] && /etc/init.d/pdnsd disable
