# target/linux/mediatek/base-files/etc/hotplug.d/gmac/00-normol
#!/bin/sh

. /lib/platform.sh

wan_def_ifname="$(uci get network.wan.def_ifname)"
wan_ifname="$(uci get network.wan.ifname)"
wan_port=4

if [ -z "$wan_def_ifname" ]; then
	wan_def_ifname="eth2.2"
fi

gmac_hotplug_cb() {
	for script_file in /etc/hotplug.d/gmac/$2/* ; do
		$script_file $1 $3
	done
}

# wisp mode skip this event
if [ "$wan_def_ifname" == "$wan_ifname" ]; then
	if [ "$PORTNUM" -eq "$wan_port" ]; then
		case "$ACTION" in
		"linkup")
			gmac_hotplug_cb start wan $wan_ifname
			;;
		"linkdown")
			gmac_hotplug_cb stop wan $wan_ifname
			;;
		esac
	fi
fi
