# customer/packages/others/opsp-client/files/opsp-client.init
#!/bin/sh /etc/rc.common
#author zk


START=97

start(){
    opsp-client -c /etc/opsp-client.conf -d &
}
stop(){
    killall opsp-client
}
restart(){
    stop
    start
}
