# target/linux/mediatek/base-files/etc/uci-defaults/hiwifi
#!/bin/sh
#
# Copyright (C) 2012-2014 www.hiwifi.com
#

. /etc/functions.sh
. /lib/platform.sh
. /lib/functions/uci-defaults.sh

hiwifi_storage_setup() {
	local board="$1"
	case $board in
	"HC5662" | "HC5661s" | "HC5761s")
		ucidef_set_hiwifi_storage "SD" "3000M"
		;;
	*)
		ucidef_set_hiwifi_storage "NULL" "0M"
		;;
	esac
}

[ ! -e /etc/config/hiwifi ] && touch /etc/config/hiwifi

board=$(tw_board_name)

hiwifi_storage_setup $board

ucidef_set_hiwifi_led_state "1"

ucidef_commit_hiwifi

exit 0
