# target/linux/mediatek/base-files/etc/uci-defaults/leds
#!/bin/sh
#
# Copyright (C) 2012-2013 hiwifi.com
#

. /lib/functions/uci-defaults.sh
. /lib/platform.sh

board=$(tw_board_name)

case $board in
HC5642)
	sed -i '/status_led_set_timer 1000 1000/d' /etc/diag.sh
	sync
	;;
esac

ucidef_set_led_default "inet" "INET" "mtk:green:internet" "0"
ucidef_set_led_netdev "wlan" "WLAN" "mtk:green:wlan-2p4" "ra0"
ucidef_set_led_netdev "wlan5g" "WLAN5G" "mtk:green:wlan-5p" "rai0"

ucidef_commit_leds

exit 0
