# target/linux/mediatek/base-files/etc/uci-defaults/network
#!/bin/sh
#
# Copyright (C) 2012-2014 www.hiwifi.com
#

. /etc/functions.sh
. /lib/platform.sh
. /lib/functions/uci-defaults.sh

platform_setup_interfaces()
{
	local board="$1"

	ucidef_set_interface_loopback

	case $board in
	HC5662 | HC5661s | HC5761s)
		ucidef_set_interfaces_lan_wan "eth2.1" "eth2.2"
		;;
	HB5801)
		ucidef_set_interface_lan_dhcp "eth2.1" "$board"
		;;
	HC5642)
		ucidef_set_interface_lan "eth2.1"
		;;
	*)
		ucidef_set_interfaces_lan_wan "eth2" "eth3"
		;;
	esac
}

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

board=$(tw_board_name)

platform_setup_interfaces $board

uci commit network

exit 0
