# target/linux/mediatek/base-files/etc/uci-defaults/wlan
#!/bin/sh
# Copyright (C) 2013-2014 www.hiwifi,com
#

. /lib/platform.sh

board=$(tw_board_name)
mac_suffix=$(tw_get_mac)
mac_suffix=${mac_suffix:6:6}
macaddr=$(ifconfig eth2 | grep HWaddr | awk '{ print $5 }' | awk -F: '{printf "D6:"$2":"$3":"$4":"$5":"$6}')

case "$board" in
"HC5642")
	sed -i "s/^MacAddress=/MacAddress=${macaddr}/g" /etc/wlan/mt7628/mt7628.dat
	sed -i 's/wpad -s/wpad -c/g' /etc/init.d/wpad
	smt=$(uci -q get wireless.radio0.smt)
	devid=$(uci -q get wireless.radio0.devid)
	uci set wireless.radio0.channel=6
	uci set wireless.radio0.repeater="1"
	[ -n "$devid" ] || uci set wireless.radio0.devid="01020200"
	[ -n "$smt" ] || {
		uci set wireless.radio0.smt="1"
		uci set wireless.@wifi-iface[0].ssid="STAR_"${mac_suffix}
	}
	uci set wireless.@wifi-iface[0].dissmt="1"
	uci set wireless.@wifi-iface[0].macaddr=${macaddr}
	uci set wireless.@wifi-iface[1].macaddr="$(ifconfig eth2 | grep HWaddr | awk '{ print $5 }')"
	;;
"HB5801")
        sed -i 's/"_5G"//g' /lib/wifi/mediatek.sh
        sed -i 's/wpad -s/wpad -c/g' /etc/init.d/wpad
        sed -i 's/^HT_BW=1/HT_BW=0/g' /etc/wlan/mt7628/mt7628.dat
        sed -i 's/^MaxStaNum=32/MaxStaNum=64/g' /etc/wlan/mt7628/mt7628.dat
        sed -i 's/^MaxStaNum=32/MaxStaNum=64/g' /etc/wlan/iNIC_ap.dat
        ;;
esac

uci rename wireless.@wifi-iface[0]=master
uci rename wireless.@wifi-iface[1]=slave
uci rename wireless.@wifi-iface[2]=guest
uci commit wireless
exit 0
